# launchd/com.mini-on-ai.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Pipeline daemon: runs continuously, pausing at Telegram approval gate -->
  <key>Label</key>
  <string>com.mini-on-ai.pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/minion/Dev/mini-on-factory/scripts/run_daemon.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/minion/Dev/mini-on-factory</string>

  <!-- Restart automatically if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Start when loaded -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/minion/Dev/mini-on-factory/logs/pipeline.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/minion/Dev/mini-on-factory/logs/pipeline-error.log</string>
</dict>
</plist>
